# etc/systemd/system/humidor-lcd.service
[Unit]
Description=Smart humidor LCD

[Service]
Type=simple
ExecStart=/usr/bin/sh-lcd
Restart=on-failure

[Install]
WantedBy=multi-user.target
